# subsonic-installer-debian/src/DEBIAN/postrm
#! /bin/sh

set -e

# Remove symlink.
rm -f /usr/bin/subsonic

# Remove startup scripts.
update-rc.d -f subsonic remove
