# libresonic-installer-debian/src/DEBIAN/prerm
#! /bin/sh

set -e

# Stop Libresonic service.
if [ -e /etc/init.d/libresonic ]; then
  invoke-rc.d libresonic stop
fi
